# debian/postinst
#!/bin/bash -ve

systemctl daemon-reload
systemctl enable gdpr-authenticator
systemctl start gdpr-authenticator
