# debian/prerm
#!/bin/bash

systemctl stop gdpr-authenticator
systemctl disable gdpr-authenticator
